# packages/ddns-openrouter/files/etc/init.d/ddns
#!/bin/sh /etc/rc.common

START=94

start () {
  sh /usr/lib/ddns/dynamic_dns_updater.sh start 
}

stop () {

  echo "ddns stop"
}

